# cstcmp_supervisor.conf
[program:cstcmp]
command = /home/pi/run/fun/cstcmp/venv/bin/gunicorn app.cstcmp:app -b localhost:8000
directory = /home/pi/run/fun/cstcmp
user = pi
